# scripts/com.myhealth.auto-beans.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.myhealth.auto-beans</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/marc.philipps/Documents/myHealth/scripts/auto-beans.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>WorkingDirectory</key>
    <string>/Users/marc.philipps/Documents/myHealth</string>

    <key>StandardOutPath</key>
    <string>/Users/marc.philipps/Documents/myHealth/logs/auto-beans-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/marc.philipps/Documents/myHealth/logs/auto-beans-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
